# docker-restic.cron
0 0 * * * docker-restic backup && docker-restic backup-prune && docker-restic backup-check && docker-restic sync