# docker-restic.cron
* * * * * echo "Docker-Restic is ready to use..."
